# conf/supervisor.conf
[program:uwsgi]
user = root
command=uwsgi -y /root/keybin/uwsgi.yml
autostart=true
autorestart=true
stderr_logfile = /var/log/uwsgi_err.log
stdout_logfile = /var/log/uwsgi_out.log
stopsignal=INT